# package/base-files/files/etc/init.d/wireless_macfilter
#!/bin/sh /etc/rc.common
USE_PROCD=1
START=97

start_service() {
	/sbin/macfilter all
}