# package/network/services/samba3/files/samba-nmbd.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2011 OpenWrt.org

START=61

start() {
	/usr/sbin/nmbd -D
}

stop() {
	killall nmbd
}
