# package/network/utils/iptables/files/firewall6.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2008-2010 OpenWrt.org

START=45

boot() {
	/sbin/firewall6 boot
}

start() {
	### applying after WAN connection up ###
	exit 0
}

stop() {
	/sbin/firewall6 stop
}

reload() {
	/sbin/firewall6 reload
}

