# package/network/utils/iptables/files/nat.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2008-2010 OpenWrt.org

START=50

start() {
	### applying after WAN connection up ###
	exit 0
}

apply() {
	/sbin/nat apply
	CONNTRACK_C=$(command -v conntrack)
	[ -n "$CONNTRACK_C" ] && $CONNTRACK_C -F
}

init() {
	exit 0
}

reload() {
	/sbin/nat reload
	CONNTRACK_C=$(command -v conntrack)
	[ -n "$CONNTRACK_C" ] && $CONNTRACK_C -F
}

